# root/etc/inittab
# Console
::respawn:-/bin/sh

# Startup
::sysinit:/etc/init.d/rcS

# Shutdown
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a
